# crontab
0 8 * * * python manage.py send_reminders
